# config_service/myhouse.service
[Unit]
Description=My house w/ raspberry pi: startup

[Service]
ExecStart=/usr/local/sbin/myhouse

[Install]
WantedBy=multi-user.target
